# cron/strava.cron
# cron/strava.cron
# 05:00 Europe/London daily
0 5 * * * /bin/sh -lc "mkdir -p /repo/logs/scheduler && /repo/scripts/flow.sh >> /repo/logs/scheduler/cron.log 2>&1"
